# es-modules-refactor.timer
[Unit]
Description=Run ES Modules refactor agent hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
