# srv.service
[Unit]
Description=FeedReader - RSS/Scraper Feed Reader

[Service]
Type=simple
User=exedev
Group=exedev
WorkingDirectory=/home/exedev/feedreader
ExecStart=/home/exedev/feedreader/feedreader
Restart=always
RestartSec=5
Environment=HOME=/home/exedev
Environment=USER=exedev
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
