# task-agent.service
[Unit]
Description=Beads task agent (shelley + claude-opus-4.6)
After=network.target beads-dolt.service

[Service]
Type=oneshot
User=exedev
WorkingDirectory=/home/exedev/feedreader
ExecStart=/home/exedev/feedreader/scripts/task-agent.sh
TimeoutStartSec=3900
Environment=HOME=/home/exedev
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/exedev/.local/bin:/home/exedev/go/bin
